# supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command=gunicorn app:app

[program:worker]
command=celery --app tasks worker --loglevel info --concurrency 4